# package/base-files/files/etc/uci-defaults/Zzzdefaults
uci set luci.main.lang=en
uci commit luci
 
uci set network.lan.ipaddr=192.168.2.1
uci set network.lan.netmask=255.255.254.0
uci commit network

uci set luci.themes.Open=/luci-static/argon		
uci	set luci.main.mediaurlbase=/luci-static/argon		
commit luci

uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci commit system

uci set system.@system[0].hostname='MasterBuild'
uci commit system

exit 0
